# build/debian-ftpbeat-amd64.service
[Unit]
Description=ftpbeat
Documentation=https://www.elastic.co/guide/en/beats/ftpbeat/current/index.html
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/share/ftpbeat/bin/ftpbeat -c /etc/ftpbeat/ftpbeat.yml -path.home /usr/share/ftpbeat -path.config /etc/ftpbeat -path.data /var/lib/ftpbeat -path.logs /var/log/ftpbeat
Restart=always

[Install]
WantedBy=multi-user.target
